# Fedora_Sway_Atomic/.config/systemd/user/backup_toolboxes.service
[Unit]
Description=Run daily Toolboxes backup

[Service]
Type=oneshot
ExecStart=/home/gerben/.local/bin/backup-toolboxes
